# supervisor/conf.d/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770


[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

loglevel=error

[program:php-fpm7]
command=/usr/sbin/php-fpm7 --nodaemonize -R

#[program:php-fpm]
#command=/usr/sbin/php-fpm --nodaemonize -R

[program:nginx]
command=/usr/sbin/nginx

[program:redis]
command=redis-server --requirepass develop
;[supervisorctl]
;serverurl=http://127.0.0.1:9723 ; use a unix:// URL  for a unix socket
[supervisorctl]              ;这个主要是针对supervisorctl的一些配置
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[include]
files = /etc/supervisor/conf.d/*.conf
